# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgres]
command=/usr/lib/postgresql/18/bin/postgres -D /data/postgres -c config_file=/data/postgres/postgresql.conf
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:minio]
command=/usr/local/bin/minio server /data/minio --console-address :9001
autostart=true
autorestart=true
priority=20
environment=MINIO_ROOT_USER="%(ENV_MINIO_ROOT_USER)s",MINIO_ROOT_PASSWORD="%(ENV_MINIO_ROOT_PASSWORD)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:app]
command=/usr/bin/node /app/apps/web/server.js
autostart=true
autorestart=true
priority=30
startsecs=5
startretries=10
environment=NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
